# dist/etc/systemd/system/awsmock.service
[Unit]
Description=AWS Simulation Service Manager
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=root
ExecStart=/usr/local/awsmock/bin/awsmockmgr --loglevel info --config /usr/local/awsmock/etc/awsmock.json
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
